# .docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -F
autorestart=true
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
priority=20

[program:queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=120
user=www-data
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:scheduler]
command=sh -c "while true; do php /var/www/html/artisan schedule:run --no-interaction; sleep 60; done"
user=www-data
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
